# infra/akash/vault-agent/supervisord.conf
; ===========================================================================
; Supervisord config for the AgentSwarm Akash container.
;
; Process tree (all unprivileged, uid 10001):
;   vault-agent     <-- writes role-id consumed, then renders templates
;   wait-for-secrets <-- one-shot, exits when env file is valid
;   agentswarm-app  <-- starts only after wait-for-secrets exits 0
;
; supervisord runs in foreground (-n in CMD) so tini reaps it cleanly.
; ===========================================================================

; supervisord itself runs as the agentswarm user (set by the Dockerfile's
; USER directive). The `user=` field is therefore intentionally omitted from
; the [supervisord] and [program:*] sections - supervisord can only setuid
; when started as root, which we explicitly do NOT do.
[supervisord]
nodaemon=true
logfile=/var/log/agentswarm/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info
pidfile=/var/run/vault-agent/supervisord.pid

; ---------------------------------------------------------------------------
; 1. Vault Agent (renders env files; signals app on rotation)
; ---------------------------------------------------------------------------
[program:vault-agent]
command=/usr/local/bin/vault agent -config=/etc/vault-agent/agent.hcl
priority=10
autostart=true
autorestart=true
startsecs=3
stopwaitsecs=10
stdout_logfile=/var/log/agentswarm/vault-agent.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/agentswarm/vault-agent.err
stderr_logfile_maxbytes=10MB
environment=VAULT_SKIP_VERIFY="false"

; ---------------------------------------------------------------------------
; 2. The actual workload. wait-for-secrets.sh is invoked INLINE so the gate
;    is a true precondition: the python process only starts if Vault Agent
;    has rendered a non-placeholder env file within the timeout.
;    AGENTSWARM_APP_CMD is set by docker-entrypoint.sh and defaults to
;    `python -m agentswarm`.
; ---------------------------------------------------------------------------
[program:agentswarm-app]
command=/bin/bash -c "/usr/local/bin/wait-for-secrets.sh && set -a && source ${AGENTSWARM_ENV_FILE} && set +a && exec ${AGENTSWARM_APP_CMD}"
directory=/opt/agentswarm
priority=20
autostart=true
autorestart=true
startsecs=10
startretries=5
stopwaitsecs=30
stopsignal=TERM
stdout_logfile=/var/log/agentswarm/app.log
stdout_logfile_maxbytes=20MB
stderr_logfile=/var/log/agentswarm/app.err
stderr_logfile_maxbytes=20MB
environment=PYTHONUNBUFFERED="1",AGENTSWARM_ENV_FILE="/etc/agentswarm/secrets/runtime.env"
